# systemd/ink.service
[Unit]
Description=Run ink, creating new backups if necessary

[Service]
Type=simple
ExecStart=+/usr/bin/env run_ink
